# ansible/templates/monitor.service.j2
[Unit]
Description=Python Monitoring Script
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/monitor.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target